# cortex.service
[Unit]
Description=Cortex - Distributed AI mesh with Raft consensus and MCP integration
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=simple
User=paschal
WorkingDirectory=/home/paschal/claudecluster
ExecStart=/usr/bin/node /home/paschal/claudecluster/dist/index.js
Restart=always
RestartSec=3
TimeoutStopSec=15
KillMode=mixed
ExecStopPost=-/usr/bin/fuser -k 50051/tcp
StandardOutput=journal
StandardError=journal

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
