# config/cortex-public-replica.service
[Unit]
Description=Refresh Cortex public replica for Litestream

[Service]
Type=oneshot
WorkingDirectory=/home/paschal/cortex
ExecStart=/usr/bin/npx tsx scripts/create-public-replica.ts
User=paschal
